# xserver/files/novnc.socket
[Unit]
Description=Go vnc socket

[Socket]
ListenStream=3050
NoDelay=true

[Install]
WantedBy = sockets.target
